# site/ovirt_jenkins_slave/templates/format-local-disk.service.erb
[Unit]
Description=Format locally attached disk
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/mkfs -t ext4 <%= @device %>
