# site/ovirt_jenkins_slave/templates/srv-local.mount.erb
[Unit]
Description=Mount locally attached disk
After=format-local-disk.service
Before=local-fs.target
Requires=format-local-disk.service
        
[Install]
WantedBy=local-fs.target
        
[Mount]
What=/dev/vdb
Where=/srv/local
